# alpineroot/etc/services.d/Xorg/run
#!/usr/bin/with-contenv bash

# Execute Xorg
/usr/bin/Xorg :1 \
	-displayfd 3 \
	-config /etc/X11/xrdp/xorg.conf \
	-noreset \
	-nolisten tcp
