# root/etc/services.d/xfce/run
#!/usr/bin/with-contenv bash

su abc -c 'DISPLAY=:0 startxfce4'
