# root/etc/services.d/vnc/run
#!/usr/bin/with-contenv bash

/usr/bin/x11vnc -shared -forever -display :0.0 -rfbport 5900
